# coco/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:start]
command=python3 /root/coco/run_server.py
startsecs=60
