# blog/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:start]
command=hexo server -i 0.0.0.0 -p 80
startsecs=60
