# luna/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:start]
command=python3 /root/luna/run_server.py
startsecs=60
